# scripts/launchd/com.stockfolio.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stockfolio.backend</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/Shared/workspace/stockfolio/backend/.venv/bin/uvicorn</string>
        <string>app.main:app</string>
        <string>--app-dir</string>
        <string>/Users/Shared/workspace/stockfolio/backend</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8000</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/Shared/workspace/stockfolio/backend</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/Shared/workspace/stockfolio/logs/backend.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/Shared/workspace/stockfolio/logs/backend.log</string>
</dict>
</plist>
